# opilio-daemon/assets/opilio-daemon.service
[Unit]
Description=opilio-daemon

[Service]
ExecStart=/usr/local/bin/opilio-daemon
Type=exec
Restart=on-failure

[Install]
WantedBy=default.target
